# Shadow/cometcloud-local/scripts/com.cometcloud.cis_scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cometcloud.cis_scheduler</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Volumes/CometCloudAI/cometcloud-local/venv/bin/python</string>
        <string>/Volumes/CometCloudAI/cometcloud-local/cis_scheduler.py</string>
        <string>--run-once</string>
        <string>--full</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Volumes/CometCloudAI/cometcloud-local/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>INTERNAL_TOKEN</key>
        <string>cometcloud_internal_2026</string>
        <key>COINGECKO_API_KEY</key>
        <string>CG-Rv47zv5eFuL2N9DXunbEQsrK</string>
        <key>RAILWAY_URL</key>
        <string>https://web-production-0cdf76.up.railway.app</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>WorkingDirectory</key>
    <string>/Volumes/CometCloudAI/cometcloud-local</string>

    <key>StandardOutPath</key>
    <string>/Volumes/CometCloudAI/cometcloud-local/_logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Volumes/CometCloudAI/cometcloud-local/_logs/launchd_stderr.log</string>
</dict>
</plist>
